# bot11.conf
[program:bot11]
directory=/home/ubuntu/bots/bot11
command=/home/ubuntu/bots/bot11/venv/bin/python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/bots/bot11/logfile_err.log
stdout_logfile=/home/ubuntu/bots/bot11/logfile.log